# deploy/proforin-bot.service
[Unit]
Description=Kliko Telegram bot worker
After=network-online.target kliko-site.service
Wants=network-online.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/srv/kliko-site/current
EnvironmentFile=/etc/kliko-site.env
ExecStart=/srv/kliko-site/venv/bin/python manage.py telegram_bot
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/srv/kliko-site
ProtectHome=true
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
